# etc/NetworkManager/dispatcher.d/22-log-connectivity
#!/bin/sh

if [ "$2" = "up" ]; then
    /usr/sbin/olpc-connectivity &
fi
