# etc/NetworkManager/dispatcher.d/22-log-connectivity
#!/bin/sh

if [ "$2" = "up" ]; then
    export PATH=$PATH:/usr/sbin
    olpc-connectivity &
fi
